# packages/init/rustysd/files/unitfiles/default.target
[Unit]
Description= Reached if boot is finished
