# packages/init/rustysd/files/unitfiles/init.target
[Unit]
Description= Base initialization for booting up the system

[Install]
WantedBy=onboot.target
