# packages/init/rustysd/files/unitfiles/network.target
[Unit]
Description= Onboot services checkpoint

[Install]
WantedBy=services.target
