# packages/init/rustysd/files/unitfiles/services.target
[Unit]
Description= Startup system services

[Install]
WantedBy=default.target
